# systemd/places-mcp-http.service
[Unit]
Description=Places MCP HTTP Wrapper (Google Places/Routes/Weather)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/projects/places-mcp
EnvironmentFile=/home/admin/projects/places-mcp/.env
Environment=HTTP_PORT=8203
Environment=HTTP_HOST=0.0.0.0
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /home/admin/projects/places-mcp/dist/http_server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
